# packaging/webfilter.service
[Unit]
Description=WebFilter Proxy (MITM filtering proxy + management UI)
Documentation=https://github.com/Yjlion/mitmproxy-web-filter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=webfilter
Group=webfilter
WorkingDirectory=/opt/webfilter
ExecStart=/opt/webfilter/webfilter run --settings config/settings.json
Restart=on-failure
RestartSec=3
# The proxy holds ~2 file descriptors per active connection (client +
# upstream); raise the limit so busy browsers don't exhaust the default 1024.
LimitNOFILE=65536
TimeoutStopSec=15

# Basic hardening - safe defaults for a network-facing daemon that only
# needs to read/write its own working directory.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/webfilter

[Install]
WantedBy=multi-user.target
